# config/deploy/shared/start.sh.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          <NAME>
# Required-Start:    $local_fs $network $named $time $syslog
# Required-Stop:     $local_fs $network $named $time $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       <DESCRIPTION>
### END INIT INFO
 
PIDFILE="<%= fetch(:pid_file) %>"

SCRIPT="<%= fetch(:env) %>/bin/python <%= fetch(:app_path) %>"
RUNAS=

LOGFILE="<%= fetch(:log_file) %>"
 
start() {
  if [ -f $PIDFILE ] && kill -0 $(cat $PIDFILE); then
    echo '<%= fetch(:application) %> already running' >&2
    return 1
  fi
  echo 'Starting <%= fetch(:application) %>' >&2
  $SCRIPT&
  echo '<%= fetch(:application) %> started' >&2
}
 
stop() {
  if [ ! -f "$PIDFILE" ] || ! kill -0 $(cat "$PIDFILE"); then
    echo '<%= fetch(:application) %> not running' >&2
    return 1
  fi
  echo 'Stopping <%= fetch(:application) %>' >&2
  kill -15 $(cat "$PIDFILE") && rm -f "$PIDFILE"
  echo '<%= fetch(:application) %> stopped' >&2
}
 
case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
esac
